# crontab
0 0 1 * * /funmoney/cron.sh >/dev/null 2>&1
